# projects/pi5-trading-system/deployment/pi5-trading.service
[Unit]
Description=Pi5 Trading System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/%u/pi5-trading-system
ExecStart=/home/%u/pi5-trading-system/api
Restart=always
RestartSec=10
StandardOutput=append:/home/%u/pi5-trading-system/logs/trading.log
StandardError=append:/home/%u/pi5-trading-system/logs/trading-error.log

# Environment
Environment="PATH=/usr/local/go/bin:/usr/bin:/bin"
EnvironmentFile=/home/%u/pi5-trading-system/.env

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=200%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
